# linux/systemd-cron/script.service
[Unit]
Description=Execute a bash script for SERVICE_NAME

[Service]
Type=oneshot
ExecStart=
